# deploy/forum-fee-reconcile.timer
[Unit]
Description=Run forum-fee-reconcile daily

[Timer]
OnBootSec=15min
OnUnitActiveSec=24h
Persistent=true

[Install]
WantedBy=timers.target
