# deploy/forum-agora-mind.service
[Unit]
Description=Forum AgoraMind keeper - AI-driven, receipt-backed, TrackRecordV2
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/forum/keeper
ExecStart=/opt/forum/keeper/node_modules/.bin/tsx /opt/forum/keeper/scripts/agora-mind-keeper.mjs --markets 3 --interval 60 --publish-every 10 --label forum-agora-mind-prod-v1 --receipts-dir /opt/forum/web/receipts --receipts-base-url https://forum.gudman.xyz/receipts
Restart=on-failure
RestartSec=30s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=forum-agora-mind

[Install]
WantedBy=multi-user.target
